# overlays/base-embedded/etc/profile.d/enforce_plasma_style.sh
export QT_QUICK_CONTROLS_STYLE=Plasma
export QT_QUICK_CONTROLS_MOBILE=1

